# packaging/deb/debian/cd-network-manager-daemon.prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop and disable the daemon service
        if [ -d /run/systemd/system ]; then
            systemctl stop cd-network-managerd.service >/dev/null || true
            systemctl disable cd-network-managerd.service >/dev/null || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
